# upHILkeyboard.service
[Unit]
Description=Create USB gadgets

[Service]
Type=oneshot
ExecStart=/share/usb_emu.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
